# hackathon/coach/install/com.claudecoach.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.claudecoach.daily</string>

	<key>ProgramArguments</key>
	<array>
		<string>__PYTHON__</string>
		<string>__REPO__/coach/scripts/coach_run.py</string>
		<string>--once</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>RunAtLoad</key>
	<false/>

	<key>StandardOutPath</key>
	<string>__HOME__/.claude/coach/run.log</string>

	<key>StandardErrorPath</key>
	<string>__HOME__/.claude/coach/run.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>__PATH__</string>
	</dict>
</dict>
</plist>
